# sdbuild/packages/usbgadget/usbgadget.service
[Unit]
Description=USB Gadget for Networking
Before=network-online.service isc-dhcp-server.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/bin/usbgadget
ExecStop=/usr/local/bin/usbgadget_stop

[Install]
WantedBy=basic.target
