# sdbuild/packages/wifi_connect/wifi_connect.service
[Unit]
Description=Wifi Autoconnection
Wants=networking.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/wifi_connect.sh

[Install]
WantedBy=basic.target
