# sdbuild/packages/clear_pl_statefile/clear_pl_statefile.service
[Unit]
Description=Clears the PL Statefile on boot
Requires=jupyter.service
After=jupyter.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/clear_pl_statefile.sh

[Install]
WantedBy=basic.target
