# sdbuild/packages/jupyter/jupyter.service
[Unit]
Description=Jupyter Notebook Server

[Service]
Type=forking
ExecStart=/usr/local/bin/start_jupyter.sh

[Install]
WantedBy=basic.target
