# sdbuild/packages/wpa_ap/wpa_ap.service
[Unit]
Description=Wifi Access Point
Before=network-online.service isc-dhcp-server.service

[Service]
Type=forking
ExecStart=/usr/local/share/wpa_ap/ap_start.sh
ExecStop=/usr/local/share/wpa_ap/ap_stop.sh

[Install]
WantedBy=basic.target
