# sdbuild/packages/bootpy/bootpy.service
[Unit]
Description=Executing boot.py from the boot partition
Requires=jupyter.service
After=jupyter.service boot.mount

[Service]
Type=oneshot
ExecStart=/usr/local/bin/bootpy.sh

[Install]
WantedBy=basic.target
