# sdbuild/packages/resizefs/resizefs.service
[Unit]
Description=Resize Filesystem on SD card

[Service]
Type=oneshot
ExecStart=/usr/local/bin/resizefs.sh

[Install]
WantedBy=basic.target
