# boards/ZCU104/packages/boot_leds_zcu104/boot_leds.service
[Unit]
Description=Boot LED flashing
Requires=pl_server.service jupyter.service
After=pl_server.service jupyter.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/boot_leds.sh

[Install]
WantedBy=basic.target
